# util/js8cli.service
[Unit]
Description=JS8CLI service

[Service]
Type=simple
ExecStart="/opt/threesixes/js8cli/js8cli daemon"

[Install]
WantedBy=multi-user.target